# git/lotns_server_setup/etc/systemd/system/soracloud.service
[Unit]
Description=SoraCloud cron.php job

[Service]
User=www-data
ExecStart=/usr/bin/php -f /var/www/soracloud/cron.php
KillMode=process